# apps/base/kaniko-archlinux/archlinux/rootfs/etc/systemd/system/getty@tty10.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin oleg %I $TERM